# examples/com.example.codex-batch-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.codex-batch-runner</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>-m</string>
    <string>codex_batch_runner</string>
    <string>--config</string>
    <string>/path/to/codex-batch-runner/examples/config.example.json</string>
    <string>run-next</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/path/to/codex-batch-runner/src</string>
  </dict>

  <key>WorkingDirectory</key>
  <string>/path/to/codex-batch-runner</string>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>StandardOutPath</key>
  <string>/path/to/codex-batch-runner/.codex-batch-runner/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/codex-batch-runner/.codex-batch-runner/launchd.err.log</string>
</dict>
</plist>
